# etc/supervisor/conf.d/process1.conf
[program:process1]
command=/process/p1.bin
autostart=true
autorestart=true
stderr_logfile=/tmp/p1e.log
stdout_logfile=/tmp/p1o.log
